# deploy/coffeemaster-web.service
[Unit]
Description=CoffeeMaster9000 Web App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/coffeemaster9000/src
Environment=PYTHONPATH=/opt/coffeemaster9000/src
EnvironmentFile=-/etc/default/coffeemaster9000
ExecStart=/opt/coffeemaster9000/.venv/bin/python /opt/coffeemaster9000/src/webapp/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
